# rpi/detector.service
[Unit]
Description=Thrifty-X (legacy) detect-from-fastcard pipeline
Documentation=https://github.com/Denny-Hwang/thrifty-x
After=network-online.target time-sync.target
Wants=network-online.target time-sync.target

[Service]
ExecStart=/home/pi/thrifty/rpi/fastdet.sh
WorkingDirectory=/home/pi/thrifty
User=pi
Group=pi
Type=simple

# Restart on transient failure (USB enumeration, brief stream stalls).
Restart=on-failure
RestartSec=120
StartLimitIntervalSec=600
StartLimitBurst=5

# Hardening — intentionally minimal.
#
# This unit's entrypoint (rpi/fastdet.sh) does two things that the
# stronger sandboxing flags applied to thriftyx-capture@.service would
# block:
#   * `sudo systemctl stop ntp` / `sudo ntpd -gq`  — requires the
#     ability to elevate, which `NoNewPrivileges=true` forbids.
#   * `mkdir -p toad/ card/ log/` under `/home/pi/detector`  — writes
#     to $HOME, which `ProtectHome=read-only` forbids.
# We therefore omit those two flags here.  ProtectSystem and
# PrivateTmp are still safe because fastdet writes only to its own
# subdirectory under $HOME and never touches /usr or /tmp.
#
# For a fully hardened deployment, prefer thriftyx-capture@.service
# instead of this legacy pipeline.
ProtectSystem=full
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=thriftyx-detector

[Install]
WantedBy=multi-user.target
